# example/rootfs/etc/s6-overlay/s6-rc.d/node/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant App: InfluxDB3 Explorer
# Runs the node daemon
# ==============================================================================

bashio::log.info "Starting Node...."

exec node /app-root/_backend_app/main.js